# deploy/systemd/polymarket-bot.service
[Unit]
Description=Polymarket Trading Bot
After=network.target

[Service]
Type=simple
User=polymarket
WorkingDirectory=/opt/polymarket-bot
EnvironmentFile=/opt/polymarket-bot/.env
ExecStart=/opt/polymarket-bot/.venv/bin/python bot_continuous.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
