# deploy/systemd/polymarket-bot.service
[Unit]
Description=Polymarket Trading Bot (Docker Compose)
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/polymarket-bot

# Remove stale containers/networks before starting (- allows graceful failure)
ExecStartPre=-/usr/bin/docker compose down --remove-orphans

# Start bot in detached mode, forcing container recreation to avoid name conflicts
ExecStart=/usr/bin/docker compose up -d --force-recreate bot

# Properly tear down containers and networks on stop
ExecStop=/usr/bin/docker compose down

# Support `systemctl reload` for quick rebuild + restart
ExecReload=/usr/bin/docker compose up -d --force-recreate --build bot

[Install]
WantedBy=multi-user.target
